# rootfs/etc/services.d/squid/run
#!/command/execlineb -P
with-contenv
# Run Squid in foreground with light verbosity so Docker/s6 capture logs.
exec /usr/sbin/squid -f /etc/squid/squid.conf -N -d1
