# rootfs/etc/services.d/squid/run
#!/command/execlineb -P
with-contenv {
  mkdir -p /var/spool/squid /var/log/squid /var/run/squid /var/logs
  touch /var/logs/cache.log /var/logs/access.log
  chown -R proxy:proxy /var/spool/squid /var/log/squid /var/run/squid /var/logs
  rm -f /var/run/squid.pid
}
exec /command/s6-setuidgid proxy /usr/sbin/squid -f /etc/squid/squid.conf -N
