# rootfs/etc/services.d/squid/run
#!/usr/bin/execlineb -P
with-contenv
  s6-setuidgid proxy
  exec squid -f /etc/squid/squid.conf -N
