# rootfs/etc/services.d/squid/run
#!/command/execlineb -P
with-contenv
exec /command/s6-setuidgid proxy /usr/sbin/squid -f /etc/squid/squid.conf -N
